# cookbooks/attero/templates/ubuntu/upstart_service.conf.erb
# supybot

description "attero service config"
author "Niko Nevala <niko@idleplaythings.com>"

setuid attero

start on runlevel [2345]
stop on runlevel [016]

respawn

exec /opt/attero/releases/current/target/start
